# deploy/soup-rendezvous@.service
[Unit]
Description=soup-rendezvous coordinator daemon (%i)
Documentation=https://github.com/8144225309/soup-rendezvous
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
Type=simple
User=root
Group=root

# --config selects the TOML, --network picks the section by the
# systemd instance name (mainnet / signet / testnet4 / etc).
ExecStart=/usr/local/bin/soup-rendezvous \
  --config=/etc/soup-rendezvous.toml \
  --network=%i \
  daemon

Environment=RUST_LOG=soup_rendezvous=info,nostr_relay_pool=warn

Restart=on-failure
RestartSec=10

SyslogIdentifier=soup-rendezvous-%i

# Same hardening as the single-network unit.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
ProtectProc=invisible
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Per-instance state dir + matching CLN socket dir.
ReadWritePaths=/var/lib/soup-rendezvous-%i /var/lib/cln-%i
ReadOnlyPaths=/etc/soup-rendezvous.toml

LimitNOFILE=4096
TasksMax=64

[Install]
WantedBy=multi-user.target
